# server/common/oursrc/execsys/execsys-binfmt.pre
#!/bin/sh
#
# execsys-binfmt: test1 
#
# chkconfig: 2345 2 99
# description: test2
#

mount -t binfmt_misc none /proc/sys/fs/binfmt_misc
echo "-1" > /proc/sys/fs/binfmt_misc/status

# START-AUTOGENERATED: DO NOT EDIT THIS SECTION, INCLUDING THIS LINE!
# This section is populated by the script upd-execsys
# END-AUTOGENERATED: DO NOT EDIT THIS SECTION, INCLUDING THIS LINE!
