# server/common/oursrc/execsys/scripts-git.socket
[Unit]
Description=Git proxy socket

[Socket]
ListenStream=9418
Accept=yes

[Install]
WantedBy=sockets.target
