# server/common/oursrc/execsys/scripts-local-smtp@.service
[Unit]
Description=Local SMTP proxy server

[Service]
ExecStart=-/usr/sbin/local-smtp-proxy
StandardInput=socket
StandardError=journal
