# server/common/oursrc/execsys/scripts-git@.service
[Unit]
Description=Git proxy server

[Service]
ExecStart=-/usr/sbin/gitproxy.pl
StandardInput=socket
User=apache
