# server/common/oursrc/execsys/scripts-svn@.service
[Unit]
Description=SVN proxy server

[Service]
ExecStart=-/usr/sbin/svnproxy.pl
StandardInput=socket
User=apache
