# server/common/oursrc/execsys/scripts-svn.socket
[Unit]
Description=SVN proxy socket

[Socket]
ListenStream=3690
Accept=yes

[Install]
WantedBy=sockets.target
