# server/common/oursrc/tokensys/crontab
@reboot	afsagent /home/afsagent/renew
0 */3 * * * afsagent /home/afsagent/renew
