# server/common/oursrc/tokensys/crontab
@reboot	/home/afsagent/renew
0  */3  *  *  *	/home/afsagent/renew
