# lvs/debian/config/etc/cron.d/restart-inetd
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
*/5 * * * * root /etc/init.d/inetd restart >/dev/null 2>&1
