# server/common/oursrc/shackle/shackle.service.in
[Unit]
Description=Shackle DBL DNS resolver socket
After=shackle.socket
Requires=shackle.socket

[Service]
ExecStart=@sbindir@/shackle 127.0.0.1 54 INET STREAM INET DGRAM
NonBlocking=true
User=nobody
Group=nobody

[Install]
WantedBy=multi-user.target
