# server/fedora/config/etc/rc.d/rc.local
#!/bin/sh

touch /var/lock/subsys/local

/bin/mkdir -m 1777 /tmp/sessions
