# server/fedora/config/etc/rc.d/rc.local
#!/bin/bash

touch /var/lock/subsys/local

/bin/mkdir -pm 1773 /tmp/sessions
