# ldap/el/config/etc/systemd/system/srv-scripts_config-ldap-el-config-etc.mount
[Unit]
Before=local-fs.target
[Mount]
What=/etc
Where=/srv/scripts_config/ldap/el/config/etc
Options=bind
[Install]
WantedBy=local-fs.target
