# server/fedora/config/etc/systemd/system/scripts-network-dependency-kludge.target
[Unit]
Requires=network.service
[Install]
WantedBy=network.target
