# server/fedora/config/etc/cron.daily/num-crontabs
#!/bin/sh

numcrontabs=$(ls -1 /var/spool/cron | wc -l)
if [ $numcrontabs -lt 100 ] && [ $numcrontabs -gt 0 ]; then
    msg="$(hostname --fqdn) has $numcrontabs files in /var/spool/cron and should have none or hundreds."
    logger -u /dev/log -p authpriv.info -t num-crontabs "$msg"
    echo "$msg"
fi
