# server/fedora/config/etc/cron.d/quickprint
MAILTO=quickprint@mit.edu
13 */6 * * * quickprint /mit/quickprint/libexec/cleanup_tmp >/dev/null
