# server/fedora/config/etc/cron.d/quickprint
MAILTO=quickprint@mit.edu
13 */6 * * * quickprint [ -e /etc/nologin ] || /mit/quickprint/bin/cleanup-tmp >/dev/null
