# server/fedora/config/etc/cron.d/scripts-cron_status
* * * * * scripts touch /afs/athena.mit.edu/contrib/scripts/cron_scripts/cron_status_flag/$(hostname -f)
