# server/fedora/config/etc/cron.d/scripts-ldap-cert-check
13 4 * * * scripts [ -x /mit/scripts/sbin/check-ldap-cert.pl ] && /mit/scripts/sbin/check-ldap-cert.pl
