# server/fedora/config/etc/cron.d/slapdagent
KRB5CCNAME=/var/run/dirsrv/krb5cc
MAILTO=scripts-root@mit.edu
0 */3 * * * fedora-ds /usr/kerberos/bin/kinit -k -t /etc/dirsrv/keytab ldap/$(hostname)
