# server/fedora/config/etc/cron.d/whoisd
@reboot root /usr/bin/twistd -y /root/whoisd/whoisd.tac
