# ansible/roles/ldirectord-status/files/ldirectord-http@.service
[Unit]
Description=Finger Per-Connection Server

[Service]
ExecStart=-/etc/ldirectord-status/ldirectord-http.sh
StandardInput=socket
