# ansible/roles/ldirectord-status/files/ldirectord-finger@.service
[Unit]
Description=Finger Per-Connection Server

[Service]
ExecStart=/etc/ldirectord-status/ldirectord-finger.sh
StandardInput=socket
