# ansible/roles/ldirectord-status/files/ldirectord-finger.socket
[Unit]
Description=Finger Socket for LDirectord Status

[Socket]
ListenStream=79
Accept=yes
TriggerLimitBurst=0

[Install]
WantedBy=sockets.target
